# config/supervisor/doi-worker.ini
[program:doi-consume]
command=php bin/console messenger:consume async_doi --time-limit=3600 --limit=10
user=pelagos
numprocs=1
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
