# config/supervisor/messenger-worker.ini
[program:messenger-consume]
command=php bin/console messenger:consume async --time-limit=3600
user=pelagos
numprocs=1
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
