# config/supervisor/messenger-worker.ini
[group:pelagos]
programs=normal-consume,low-consume,doi-consume

[program:normal-consume]
directory=/opt/pelagos
command=php bin/console messenger:consume async_normal --time-limit=3600
user=pelagos
numprocs=10
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d

[program:low-consume]
directory=/opt/pelagos
command=php bin/console messenger:consume async_low --time-limit=3600
user=pelagos
numprocs=5
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d

[program:doi-consume]
directory=/opt/pelagos
command=php bin/console messenger:consume async_doi --time-limit=3600 --limit=10
user=pelagos
numprocs=1
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
